# ops/deep-ai-research-tag-engagements.service
[Unit]
Description=deep-ai-research — tag author engagements (cheap; no external API calls)
Documentation=https://github.com/jamcas14/deep-ai-research
After=deep-ai-research-ingest.service

[Service]
Type=oneshot
WorkingDirectory=%h/code/projects/deep-ai-research
ExecStart=%h/.local/bin/uv run python -m ingest.tag_engagements -v
StandardOutput=journal
StandardError=journal
CPUWeight=20
IOWeight=20
Nice=15
RemainAfterExit=no

[Install]
WantedBy=default.target
