# ops/deep-ai-research-tag-engagements.timer
[Unit]
Description=deep-ai-research — schedule author-engagement tagging hourly

[Timer]
# After ingest catches up; cheap so we can run often.
OnBootSec=15min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
