# ops/deep-ai-research-promote-arxiv.service
[Unit]
Description=deep-ai-research — daily promoted arXiv pipeline (full-text persistence for cross-mentioned papers)
Documentation=https://github.com/jamcas14/deep-ai-research
After=network-online.target deep-ai-research-ingest.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/code/projects/deep-ai-research
ExecStart=%h/.local/bin/uv run python -m ingest.promote_arxiv -v
StandardOutput=journal
StandardError=journal
CPUWeight=30
IOWeight=30
Nice=15
RemainAfterExit=no

[Install]
WantedBy=default.target
