# ops/deep-ai-research-embed.service
[Unit]
Description=deep-ai-research — embed any unembedded chunks into sqlite-vec
Documentation=https://github.com/jamcas14/deep-ai-research
After=network-online.target deep-ai-research-ingest.service

[Service]
Type=oneshot
WorkingDirectory=%h/code/projects/deep-ai-research
ExecStart=%h/.local/bin/uv run --extra embed python -m ingest.embed_pending
StandardOutput=journal
StandardError=journal
CPUWeight=20
IOWeight=20
Nice=18
RemainAfterExit=no

[Install]
WantedBy=default.target
